# docker-service/templates/docker-compose-service
{% from "monitoring/map.jinja" import monitoring with context %}
#!/bin/bash
### BEGIN INIT INFO
# Provides:          docker-compose-{{service_name}}
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 6
# Short-Description: docker-compose {{service_name}}
# Description:       Start/stop {{service_name}} containers
### END INIT INFO

start()
{
    /usr/local/bin/docker-compose -p {{service_name}} -f /etc/docker-compose/{{service_name}}/docker-compose.yml pull
    COMPOSE_HTTP_TIMEOUT=120 /usr/local/bin/docker-compose -p {{service_name}} -f /etc/docker-compose/{{service_name}}/docker-compose.yml up --remove-orphans  -d
    ret=$?
    /usr/sbin/conntrack -D -p udp
    exit ${ret}
}

stop()
{
    /usr/local/bin/docker-compose -p {{service_name}} -f /etc/docker-compose/{{service_name}}/docker-compose.yml stop -t 30
    /usr/local/bin/docker-compose -p {{service_name}} -f /etc/docker-compose/{{service_name}}/docker-compose.yml kill
}

docker_compose() {
  /usr/local/bin/docker-compose -p {{service_name}} -f /etc/docker-compose/{{service_name}}/docker-compose.yml "$1"
}

docker_compose_status() {
  local docker_status=$(docker_compose "ps" | grep "Up" &>/dev/null)$?

  if [[ "$docker_status" == 0 ]];then
    echo "{{service_name}} is running"
  else
    docker_compose "ps"
  fi
}

case "$1" in
  start)
    start
    ;;
  stop)
    stop
    ;;
  reload)
    start
    ;;
  restart)
    docker_compose "restart"
    ;;
  status)
    docker_compose_status
    ;;
  ps)
    docker_compose "ps"
    ;;
  *)
    echo "Usage: $0 {start|stop|reload|restart|ps}"
    exit 1
    ;;
esac
